# etc/supervisord-ssh.conf
[supervisord]
nodaemon=true

[program:mongodb]
command=/usr/bin/mongod
autorestart=true

