# etc/supervisord-serf.conf
[supervisord]
nodaemon=true

[program:serf-agent]
command=/usr/bin/serf agent -config-file /docker-serf/etc/haproxy.json 
autorestart=true
